# backend/deploy/litchai-worker.service
# LitchAI Procrastinate worker — install: sudo cp to /etc/systemd/system/, daemon-reload, enable --now
# Queue state lives in Postgres (survives restarts); one worker on the 2-OCPU box.
[Unit]
Description=LitchAI Procrastinate worker
After=network.target docker.service

[Service]
User=ubuntu
# Same ciphertext store as the API — the worker reads envelopes from it.
StateDirectory=litchai
StateDirectoryMode=0700
EnvironmentFile=/etc/litchai/env
WorkingDirectory=/home/ubuntu/litchconsulting/backend
ExecStart=/home/ubuntu/litchai/venv/bin/procrastinate --app=litchai.queue.queue worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
